# systemd/mcp-proxy.service
[Unit]
Description=MCP Proxy — expõe codebase-rag via HTTP/SSE para clients remotos
After=network.target codebase-rag-server.service
Requires=codebase-rag-server.service

[Service]
Type=simple
User=gabrielsb
ExecStart=/home/gabrielsb/.local/bin/mcp-proxy --port 8080 -- /home/gabrielsb/.local/bin/uv run --directory /home/gabrielsb/mcp-servers/codebase-rag python main.py
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
